# Ingress/examples/bench/ab/build_index.cron
* * * * * /build_index.sh /siteroot/fs >> /var/log/cron.log 2>&1
